# zabbix-proxy/rootfs/etc/s6-overlay/s6-rc.d/zabbix-proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Grafana
# Runs the Grafana Server
# ==============================================================================

bashio::log.info 'Starting Zabbix-proxy...'

options+=(--foreground)
options+=(-c /etc/zabbix/zabbix_proxy.conf)

exec /usr/sbin/zabbix_proxy "${options[@]}"